# service/bluetooth-detection-home-assistant.service
[Unit]
Description=Bluetooth Detection for Home Assistant
After=network.target bluetooth.target

[Service]
Type=simple
User=root

# Change this to your actual development directory
WorkingDirectory=/home/wadda/ha_bluetooth_detect

# Use the Python from your venv
ExecStart=python /home/wadda/ha_bluetooth_detect/main.py

# Restart policy
Restart=always
RestartSec=10

# Environment variables (if you need to set any explicitly)
Environment="PATH=/usr/bin:/usr/local/bin"

# Logging
StandardOutput=journal
StandardError=journal

# Security options (optional but recommended)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
